# dot.config/profile.d/exports.sh
export EDITOR=kak
export GOPATH=$HOME/src/go
export PATH=$HOME/.local/bin:$GOPATH/bin:$PATH
export MANPATH=$HOME/.local/share/man:/usr/share/man:/usr/local/share/man
